# docker/motley_cue/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd_error.log

[program:motley_cue]
command=%(ENV_VENV_PATH)s/bin/motley_cue_uvicorn
autostart=true
autorestart=true
stdout_logfile=/var/log/motley_cue.log
stderr_logfile=/var/log/motley_cue_error.log
